# dist/slc-server/src/main/rpm/etc/init.d/slc-server
#!/bin/bash
#
# slc-server	SLC Server
#
# chkconfig: 3 92 92
# description: SLC Server

# Source function library.
. /etc/rc.d/init.d/functions

RETVAL=0

start() {
	/usr/sbin/osgi-service slc-server start
	action $"Start SLC server" /bin/true
	return $RETVAL
}

stop() {
	/usr/sbin/osgi-service slc-server stop
	action $"Stop SLC server" /bin/true
	return $RETVAL
}

# See how we were called.
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart)
        stop
        start
        RETVAL=$?
        ;;
  status)
  		/usr/sbin/osgi-service slc-server status
        ;;
  *)
        echo $"Usage: $0 {start|stop|restart|status}"
        exit 1
esac

exit $RETVAL
