# dist/slc/src/main/scripts/postinstall
if [ $1 = "1" ];then
	APP=%{name}
	/sbin/chkconfig --add $APP
	/sbin/chkconfig --level 3 $APP on
fi