# dist/slc/src/main/scripts/preremove
if [ $1 = "0" ];then
	#APP=%{name} => doesn't work
	APP=slc
	/sbin/chkconfig --del $APP
	userdel $APP
	
	# necessary for EL5? (breaks EL6)
	#groupdel slc
fi
