# dist/slc-server/src/main/scripts/postinstall
if [ $1 = "1" ];then
	APP=%{name}
	mkdir -p /var/lib/$APP/{conf,data}
	/sbin/chkconfig --add $APP
fi