# dist/slc/rpm/usr/lib/systemd/system/slcd.service
[Unit]
Description=Argeo SLC daemon
After=network.target

[Service]
ExecStart=/usr/sbin/slcd start
ExecStop=/usr/sbin/slcd stop
Type=forking
PIDFile=/var/run/slcd.pid

[Install]
WantedBy=default.target